# django_deploy_cli/services/templates/gunicorn.socket
[Unit]
Description={{ description }}

[Socket]
ListenStream=/run/{{ socket_name }}

[Install]
WantedBy=sockets.target